# dataone-cn-mnsynchronization/DEBIAN/postinst
#!/bin/bash

LONG_DATE=`date +%Y%m%d%H%M%S`

chmod +x /etc/init.d/mn-synchronize
update-rc.d mn-synchronize defaults 99
## Stop apache
echo "Start Sychronization with /etc/init.d/mn-synchronize start"
## should only be run after metacat is running and configured
##/etc/init.d/mn-synchronize start





