# dataone-cn-mnsynchronization/DEBIAN/postinst
#!/bin/bash

LONG_DATE=`date +%Y%m%d%H%M%S`

chmod +x /etc/init.d/d1-synchronize
#chmod +x /etc/init.d/d1-indexer

update-rc.d d1-synchronize defaults 99
#update-rc.d d1-indexer defaults 99

mkdir  /var/log/dataone
chown tomcat6.tomcat6 /var/log/dataone
echo "Start Sychronization with /etc/init.d/d1-synchronize start"
echo "Start Indexer with /etc/init.d/d1-indexer start"





